# etc/service/nginx/run
#!/bin/sh

# Make sure syslog is running before we start
sv start syslog-ng || exit
sleep 2

chown -R www-data:www-data /usr/share/nginx/www
exec /usr/sbin/nginx -c /etc/nginx/nginx.conf
